# rootfs/etc/s6-overlay/s6-rc.d/solarhq-scheduler/run
#!/command/with-contenv bash
exec /usr/bin/solarhq-scheduler
